# src/impl/logger.service
[Unit]
Description=Logger for MIDF services

[Service]
ExecStart=/usr/bin/logger-service
Restart=always

[Install]
WantedBy=multi-user.target
